# jobs/crontab
# Pocketcasts - every 2 hours
0 */2 * * * python3 /scripts/pocketcasts_sync.py
